# assets/etc/supervisord/dokubackup.conf
[program:dokubackup]
command = /usr/local/bin/dokubackup
user = root
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
